# humnerlenis/canaima-semilla/canaima-semilla/profiles/popular/IMG_INCLUDES/live/data/autostart/disabler.desktop
[Desktop Entry]
Type=Application
Name=Disabler
Exec=/bin/bash /lib/live/mount/medium/live/data/scripts/disabler.sh
